# roles/backup/templates/pgbackrest-follow-primary.timer.j2
[Unit]
Description=Run pgBackRest follow-primary periodically

[Timer]
OnBootSec=30
OnUnitActiveSec=30
AccuracySec=5
Unit=pgbackrest-follow-primary.service

[Install]
WantedBy=timers.target

